# systemd/tooru@.timer
[Unit]
Description=Archive and sync tooru content in %I

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
